# lnxsmartinfo19.cron
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:
#
#
# This file should be moved to /etc/cron.d/lnxsmartinfo19
# NOTE: no extension!
#
#
# m     h       dom     mon     dow     command
*/5     *       *       *       *       root  ( /home/beheer/lnxdiagd/lnxsmartinfo19.sh |logger -p info -t lnxsmartinfo19 )
@reboot                                 root  ( /home/beheer/lnxdiagd/lnxsmartinfo19.sh |logger -p info -t lnxsmartinfo19 )
